# rootfs/etc/s6-overlay/s6-rc.d/init-env/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: InvoiceNinja
# Runs the init
# ==============================================================================

# Log the start of the script
bashio::log.info "Init..."

chown www-data:www-data /var/www/app/

# Ensure necessary services start
bashio::log.info "Starting services..."